# twilio/matt_voip/voip.conf
[program:voip]
directory=/home/matt/voip
command=python2.7 voip.py
user=matt
group=matt
redirect_stderr=true
stdout_logfile=/home/matt/voip/voip.log
stdout_logfile_maxbytes=10MB
